# var/config-static/etc/cron.d/linuxmuster-backup
# /etc/cron.d/linuxmuster-backup
# Zeiten, zu denen ein Backup angefertigt wird.
#
# Frank Schütte <fschuett@gymnasium-himmelsthuer.de>
# 17.10.2013
# Gpl v3
#
0 1 * * 6   root [ -e /etc/default/linuxmuster-backup ] && . /etc/default/linuxmuster-backup && [ "$START_LINUXMUSTER_BACKUP" == "yes" ] && /usr/sbin/linuxmuster-backup-diff-full
0 2 * * 0-5 root [ -e /etc/default/linuxmuster-backup ] && . /etc/default/linuxmuster-backup && [ "$START_LINUXMUSTER_BACKUP" == "yes" ] && /usr/sbin/linuxmuster-backup --inc
